# ops/systemd/sam-ui-manager.service
[Unit]
Description=SAM UI Manager
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=samui
Group=samui
WorkingDirectory=/opt/sam-ui-manager
EnvironmentFile=-/etc/default/sam-ui-manager
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
